# ansible/files/postgresql.service.j2
# This is an additive override for the service files generated by the postgresql
# installation. Systemd will pick this directive up and append it to the default
# service definitions /lib/systemd/system/postgresql*
[Service]
Slice=slices.service
