# docker/all-in-one/salt-init/salt/supervisor/templates/gotrue.conf.j2
[program:gotrue]
directory=/opt/gotrue
command=/opt/gotrue/gotrue --config /etc/gotrue.env
user=gotrue
startretries=30
autorestart={{ gotrue_autostart }}
autostart={{ gotrue_autostart }}
stdout_logfile=/var/log/services/gotrue.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
